# config/supervisor.d/crond.ini
[program:crond]
command = crond -f
stdout_logfile=/var/log/supervisord-crond.log
stderr_logfile=/var/log/supervisord-crond.log
autostart=true
autorestart=true
