# config/supervisor.d/php-fpm.ini
[program:php-fpm]
command=/usr/sbin/php-fpm7 -F
stdout_logfile=/var/log/supervisord-php-fpm.log
stderr_logfile=/var/log/supervisord-php-fpm.log
autostart=true
autorestart=true
priority=5
stdout_events_enabled=true
stderr_events_enabled=true
